# docker/rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/command/with-contenv sh
sleep 3
cd /app
./rustdesk-console -c ./conf/config.yaml
